# systemd/streamdisplay.service
[Unit]
Description=Stream Display Server
After=network-online.target mosquitto.service
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/streamdisplay

# Keine Display-Umgebung nötig - direkter DRM/Framebuffer Zugriff
ExecStart=/opt/streamdisplay/venv/bin/python /opt/streamdisplay/app.py
Restart=always
RestartSec=5
StandardOutput=append:/var/log/streamdisplay/server.log
StandardError=append:/var/log/streamdisplay/error.log

# Zugriff auf Video-Hardware
SupplementaryGroups=video render

[Install]
WantedBy=multi-user.target
